# pvlib_benchmarks.timer
[Unit]
Description=pvlib benchmark timer

[Timer]
# every day at midnight UTC
OnCalendar=*-*-* 00:00:00

